# example/crontab
# min   hour    day     month   weekday command
*       *       *       *       *       php /srv/src/scripts/du.php
*/5     *       *       *       *       php /srv/src/scripts/gc.php
